# archive/after-install
# domain provision
samba-tool domain provision --use-rfc2307 --realm="$REALM" --domain="$DOMAIN" --server-role=dc --adminpass="$PASSWORD"

# chrony, ntpdate install and config
apt install -y chrony ntpdate
ntpdate -bu pool.ntp.org
echo "allow 0.0.0.0/0" >> /etc/chrony/chrony.conf
echo "ntpsigndsocket  /var/lib/samba/ntp_signd" >> /etc/chrony/chrony.conf
chown root:_chrony /var/lib/samba/ntp_signd/
chmod 750 /var/lib/samba/ntp_signd/
systemctl restart chrony
systemctl enable chrony

# bind9 install and config
apt install -y bind9
# named.conf.options
cat /etc/bind/named.conf.options << EOF
options {
 directory "/var/cache/bind";

 forwarders {
        8.8.8.8;
  };

 allow-query {  any;};
 dnssec-validation no;

 auth-nxdomain no;
 listen-on-v6 { any; };

 tkey-gssapi-keytab "/var/lib/samba/bind-dns/dns.keytab";
 minimal-responses yes;
};
EOF
